# cron.d/jobperf-worker
#!/bin/bash

*/5 * * * * pl543  /home/pl543/project/jobperf/scripts/worker.sh
